# workstation/macOS/home/Library/LaunchAgents/com.scry.openclaw.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.scry.openclaw.backup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>/Users/sawyer/github/scry-home/scripts/ops/daily-openclaw-backup.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/sawyer/github/scry-home</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>2</integer>
		<key>Minute</key>
		<integer>20</integer>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/sawyer/Library/Logs/scry/openclaw-backup.launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/sawyer/Library/Logs/scry/openclaw-backup.launchd.err.log</string>
</dict>
</plist>
