# workstation/macOS/home/Library/LaunchAgents/com.sawyer.scry.backups.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sawyer.scry.backups</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sawyer/github/configs/scripts/run-automated-backups.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>StandardOutPath</key>
	<string>/Users/sawyer/Library/Logs/com.sawyer.scry.backups.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/sawyer/Library/Logs/com.sawyer.scry.backups.err.log</string>
</dict>
</plist>
